# deploy/kylin/systemd/assistance-management-system.service
[Unit]
Description=帮扶管理信息系统
Documentation=file:///opt/assistance-management-system/README.md
After=network.target
Wants=network.target

[Service]
Type=simple
User=assistance-management-system
Group=assistance-management-system
WorkingDirectory=/opt/assistance-management-system

# 环境变量
EnvironmentFile=/opt/assistance-management-system/config/kylin.env

# 启动命令
ExecStart=/opt/assistance-management-system/backend/assistance-management-backend
ExecReload=/bin/kill -HUP $MAINPID

# 进程管理
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3
TimeoutStartSec=60
TimeoutStopSec=15

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=assistance-management-system

# 安全加固
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/assistance-management-system
ReadWritePaths=/var/log/assistance-management-system

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

# OOM 保护（避免被系统误杀）
OOMScoreAdjust=-500

[Install]
WantedBy=multi-user.target
